# solr/templates/solr.conf.erb
# This is the Ubuntu way to run Solr as a daemon.
# Copy to /etc/init/ to run as daemon
# Run with : sudo initctl start solr

start on runlevel [2345]
stop on runlevel [06]
respawn

exec su - <%= user %> -c 'cd /opt/solr/; java -server -Xmx1024m -Dsolr.solr.home=/etc/solr -Djetty.logs=/var/log/solr -Dsolr.data.dir=/var/lib/solr -jar start.jar >> /var/log/solr/output.log 2>&1'
